# docker/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisord.log
loglevel=info

[program:php-fpm]
command=/usr/sbin/php-fpm
autostart=true
autorestart=true
stdout_logfile=/var/log/php-fpm.log
stderr_logfile=/var/log/php-fpm.err

[program:queue-worker]
command=/usr/bin/php /var/www/html/artisan queue:work redis --sleep=3 --tries=3 --timeout=90
autostart=true
autorestart=true
stdout_logfile=/var/log/queue-worker.log
stderr_logfile=/var/log/queue-worker.err

[program:websockets]
command=/usr/bin/php /var/www/html/artisan websockets:serve
autostart=false
autorestart=true
stdout_logfile=/var/log/websockets.log
stderr_logfile=/var/log/websockets.err
